# dot_config/systemd/user/chameleos.service
[Unit]
Description=Chameleos
After=graphical-session.target

[Service]
ExecStart=/usr/bin/chameleos
Restart=always
RestartSec=2
KillMode=process
